# dot_config/systemd/user/awww-manager-random.service
[Unit]
Description=Randomize wallpapers via awww-manager
After=graphical-session.target awww.service
Wants=awww.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/awww-manager.py random
