# .config/autostart/jellyfin-mpv-shim.desktop
[Desktop Entry]
Exec=/usr/bin/jellyfin-mpv-shim
Name=jellyfin-mpv-shim
Type=Application
